# mikrotik-huawei-toolkit/mikrotik-huawei-toolkit/deploy/systemd/mkhw-backup.timer
[Unit]
Description=Executa mkhw-toolkit (backup + diff + auditoria) diariamente às 02:00

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
